# ty/jamsend-vnc.service
[Unit]
Description=x11vnc serving the real jamsend desktop (:0) — for manual FSA grants
# Replaces the old three-part HEADLESS stack (jamsend-xvfb + jamsend-wm + jamsend-x11vnc
#  on :10). Serving the REAL console desktop is what makes the native "open share" file
#   picker actually work: a real DE brings the dbus session + xdg-desktop-portal backend +
#    window focus that the bare Xvfb :10 never had (there the picker just flickered and
#     threw "The user aborted a request" / AbortError). Less to run, and it grants FSA.
# Wait for the display manager so :0 exists before we attach — the display manager ONLY.
#  After=graphical.target used to ride along and wedged a host whose boot transaction never
#   completes (a plymouth-quit-wait-style straggler holds the target's job open forever):
#    every start of this unit parked as a `waiting` job that never ran, journal empty,
#     while gdm + Xorg were visibly fine.  gdm up IS the condition we need; the target
#      added nothing but that trap.
After=display-manager.service
Wants=display-manager.service

[Service]
# Runs as root so -auth guess can read the display-manager's cookie for :0.
# Performance (the first cut's -wait 80 -defer 20 was ≈12 fps and felt sluggish on LAN):
#  -wait 10 -defer 1      poll up to 100x/s, coalesce only 1ms → snappy cursor/typing/scroll
#  -threads               threaded libvncserver — smoother client I/O
#  -ncache 10 -ncache_cr  client-side pixel cache (window moves/scroll reuse cached tiles). Costs
#                          a tall off-screen scratch strip the viewer draws BELOW the desktop —
#                          harmless, just scroll past it. Drop these two if that bugs you.
#  XDAMAGE stays ON (default): only CHANGED regions are read — essential while -noshm is in play.
# -noshm: no MIT shared memory. Attaching an SHM segment to :0 across the root/owner boundary is
#  denied (BadAccess on X_ShmAttach). This is the RELIABLE-but-slower choice; the FASTER choice is
#   to run as the session owner so SHM works — see the commented alternative below.
# -localhost: bind loopback ONLY, so the SSH tunnel (ty/desktop-viewer.sh) is the sole way in.
ExecStart=/usr/bin/x11vnc -display :0 -auth /run/user/1000/gdm/Xauthority -noshm -localhost -nopw -forever -shared -threads -ncache 10 -ncache_cr -wait 10 -defer 1 -rfbport 5900
# ^ PROVEN on this box (gdm session): root + the explicit gdm cookie + -noshm. -auth guess does NOT
#   work for :0 here. Cookie path is host/DM-specific (gdm → /run/user/<uid>/gdm/Xauthority; lightdm
#    → /var/run/lightdm/root/:0). Find it from the running Xorg's -auth arg if it differs.
#
# FASTER alternative — that gdm cookie is owned by uid 1000, so run x11vnc AS the owner and MIT-SHM
#  works: drop -noshm and each framebuffer read becomes a memcpy (the single biggest speed win). Swap
#   the two lines above for these; if x11vnc then can't attach :0, revert to the proven line.
#   User=1000
#   ExecStart=/usr/bin/x11vnc -display :0 -auth /run/user/1000/gdm/Xauthority -localhost -nopw -forever -shared -threads -ncache 10 -ncache_cr -wait 10 -defer 1 -rfbport 5900
Restart=always
RestartSec=5

[Install]
WantedBy=graphical.target
