# ty/jamsend-vnc-forward.service
[Unit]
Description=Forward VNC from jamsend VM (192.168.122.100:5910) to host 127.0.0.1:5910
Documentation=Keeps ty/xvfb-viewer.sh working unchanged now that x11vnc runs inside the VM
After=libvirtd.service jamsend-virtreset.service

[Service]
# bind=127.0.0.1: only SSH tunnels can reach this, matching x11vnc's -localhost flag inside the VM
ExecStart=/usr/bin/socat \
    TCP-LISTEN:5910,bind=127.0.0.1,reuseaddr,fork \
    TCP:192.168.122.100:5910
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
