# ty/jamsend-supervisor.service
[Unit]
Description=jamsend Watchdog and Launcher
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/s/src/jamsend/ty
# Ensure we have the display environment
Environment=DISPLAY=:0
# Start the launcher in the background, then the compose
ExecStartPre=/usr/bin/python3 /home/s/src/jamsend/ty/launcher.py &
ExecStart=/usr/bin/docker compose up
ExecStop=/usr/bin/docker compose down
ExecStopPost=/usr/bin/pkill -f launcher.py

[Install]
WantedBy=multi-user.target